# gh.service

[Service]
ExecStart=/bin/triger.bash

[Install]
WantedBy=default.target
